# deploy/pi-setup/files/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin cdreborn --noclear --skip-login --login-program /usr/local/bin/tty-login.sh %I $TERM
Type=idle
