# system/systemd/headunit-update-monitor.path
[Unit]
Description=Monitor /data/incoming_updates for new packages

[Path]
PathExistsGlob=/data/incoming_updates/*.tar.gz
Unit=headunit-update-monitor.service

[Install]
WantedBy=multi-user.target
