# deploy/pi-setup/files/etc/systemd/system/sync-before-shutdown.service
[Unit]
Description=Flush filesystem buffers before shutdown/reboot
Documentation=man:sync(1)
# Идём очень рано в последовательности остановки
DefaultDependencies=no
# Запустись до размонтирования ФС и финального shutdown
Before=umount.target shutdown.target
# Убедимся, что нужные ФС смонтированы (если есть)
RequiresMountsFor=/ /data /boot

[Service]
Type=oneshot
# Понижаем «окно риска»: временно ужесточаем commit на /data, затем двойной sync
ExecStart=/bin/sh -c 'mount -o remount,commit=1 /data 2>/dev/null || true; sync; sleep 1; sync'

[Install]
WantedBy=halt.target
WantedBy=reboot.target
WantedBy=poweroff.target
WantedBy=kexec.target
